# conf/crontab
MAILTO=$ADMIN_EMAIL
26 2 */7 * * cronic certbot renew
