# conf/crontab
MAILTO=$ADMIN_EMAIL
@reboot /srv/bin/init-run-dir.sh
# Enabled by install.sh when SSL certificates are managed by certbot:
#26 2 */7 * * cronic certbot renew
